# telco-vecchio-init
#!/bin/sh /etc/rc.common
START=99

PROG=/usr/bin/telco-vecchio

start() {
	$PROG --daemon
}

stop() {
	return 0
}
